# systemd/prometheus-co2monitor-exporter@.service
[Unit]
Description=co2monitor_exporter

[Service]
ExecStart=/usr/local/bin/co2monitor_exporter --device /dev/%I
Restart=always
RestartSec=2
